# update.cron
* * * * * python3.11 /hello_world.py >> /out.txt
* * * * * (sleep 30; python3.11 /hello_world.py >> /out.txt) 

