# process/update.cron
* * * * * python3.11 /app/hello_world.py >> /out.txt
* * * * * (sleep 30; python3.11 /app/hello_world.py >> /out.txt)
